# job diario/deploy/job-diario.service
[Unit]
Description=Job diario - Boletin Oficial Mendoza (descarga + extraccion + carga incremental)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
# Ajustar esta ruta al lugar real donde se clonó el repo en la instancia.
WorkingDirectory=/opt/boletin-oficial
EnvironmentFile=/opt/boletin-oficial/job diario/deploy/.env
ExecStart=/opt/boletin-oficial/.venv/bin/python "/opt/boletin-oficial/job diario/run_diario.py"
# Si por lo que sea la corrida se cuelga (red, API lenta), no dejar el
# proceso corriendo para siempre — cortar y que se note en journalctl.
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
